# matetherebot.timer
[Unit]
Description=Update matetherebot

[Timer]
OnBootSec=5min
OnUnitActiveSec=4m

[Install]
WantedBy=timers.target
